# configs/bt-agent.service
[Unit]
Description=Bluetooth Agent Service
After=bluetooth.service
Wants=bluetooth.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/bt-agent.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target