# configs/bt-audio-autostart.service
[Unit]
Description=Bluetooth Audio Auto-start
After=bluetooth.service network.target
Wants=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/local/bin/start-bt-audio.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target